# systemd/system/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Description=Root device
Before=local-fs.target

[Mount]
What=/dev/sda7
Where=/
Type=ext4
Options=rw,relatime,data=ordered
