# systemd/system/media-arch.automount
# Automatically generated by systemd-fstab-generator

[Unit]
Description=Automount /media/arch
Before=local-fs.target
[Automount]
Where=/media/arch
