# systemd/system/media-arch.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Description=Mount /dev/sda6 as /media/arch

[Mount]
What=/dev/sda6
Where=/media/arch
Type=ext4
Options=noauto,x-systemd.automount,user,exec,suid
