# systemd/user/user-programs.target
[Unit]
Description  = User Programs
Wants        = wm.target
After        = wm.target
AllowIsolate = true

[Install]
Alias = default.target
